# deploy/laundry-app.service
# systemd unit template for the Laundry Recipe PWA (Docker Compose).
# Placeholders (__WORKDIR__, __COMPOSE__) are filled in by deploy/install.sh.
# Do not install this file directly; run deploy/install.sh instead.

[Unit]
Description=Laundry Recipe PWA (Docker Compose)
Documentation=https://github.com/Woz4tetra/laundry-app
Requires=docker.service
After=docker.service network-online.target
Wants=network-online.target

[Service]
Type=oneshot
RemainAfterExit=yes
WorkingDirectory=__WORKDIR__

# Bring the stack up at boot (image is built ahead of time by install.sh).
# Compose's own `restart: unless-stopped` keeps the container alive after that.
ExecStart=__COMPOSE__ up -d
ExecStop=__COMPOSE__ down
# `systemctl reload laundry-app` rebuilds and redeploys after a code change.
ExecReload=__COMPOSE__ up -d --build

# `up`/`down` can take a while when images need building or pulling.
TimeoutStartSec=0

[Install]
WantedBy=multi-user.target
